# com.paranoidjarvis.startupcam.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.paranoidjarvis.startupcam</string>
    <key>ProgramArguments</key>
    <array>
    <string>/Users/exempt/Desktop/Projects/Paranoid-Jarvis/.venv/bin/python</string>
        <string>/Users/exempt/Desktop/Projects/Paranoid-Jarvis/startup_cam_logger.py</string>
    </array>
    <key>RunAtLoad</key>
    <false/>
    <key>StandardOutPath</key>
    <string>/tmp/startupcam.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/startupcam.err</string>
    <key>WorkingDirectory</key>
    <string>/Users/exempt/Desktop/Projects/Paranoid-Jarvis</string>
</dict>
</plist>
